# packaging/stoandl.service
[Unit]
Description=Stoandl — Pebble companion daemon for Linux
After=bluetooth.target network-online.target
Wants=bluetooth.target

[Service]
Type=simple
# Requires a JDK 25 runtime (libpebble3's BT Classic transport uses java.lang.foreign).
# --enable-native-access=ALL-UNNAMED silences JDK 25's restricted-native-access warning; never
# --enable-preview (FFM is final on 22+, that flag would fail on a release JDK).
ExecStart=/usr/bin/java --enable-native-access=ALL-UNNAMED -jar /usr/lib/stoandl/stoandl.jar
Restart=on-failure
RestartSec=5
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus

[Install]
WantedBy=default.target
